# scripts/ops/cc-inverter-rawjson-prune.timer
[Unit]
Description=Daily prune of inverter_readings.raw_json

[Timer]
# Run after the nightly logical backup (02:39 UTC) has typically finished.
OnCalendar=*-*-* 04:30:00 UTC
RandomizedDelaySec=300
# Persistent=false on purpose: a catch-up run on enable could collide with the
# nightly backup's table locks. Just schedule the next 04:30 slot.
Persistent=false

[Install]
WantedBy=timers.target
