# scripts/ops/systemd/cc-postgres-restore-verify.service
[Unit]
Description=1PWR CC restore verification drill
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
EnvironmentFile=/etc/default/cc-postgres-backup
ExecStart=/usr/local/bin/cc_postgres_restore_verify.sh
