# acdb-api/systemd/sms-reconcile.timer
[Unit]
Description=Run SMS gateway reconciliation every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
